# com.stockbot.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.stockbot</string>
    
    <key>ProgramArguments</key>
    <array>
        <string>/usr/bin/python3</string>
        <string>/Users/kuma/Desktop/Python/Stock_Bot/Stock_Bot_V1_main.py</string>
        <string>--run-now</string>
    </array>
    
    <key>WorkingDirectory</key>
    <string>/Users/kuma/Desktop/Python/Stock_Bot</string>
    
    <key>StartCalendarInterval</key>
    <array>
        <!-- Monday 9:15 AM -->
        <dict>
            <key>Weekday</key>
            <integer>1</integer>
            <key>Hour</key>
            <integer>9</integer>
            <key>Minute</key>
            <integer>15</integer>
        </dict>
        <!-- Tuesday 9:15 AM -->
        <dict>
            <key>Weekday</key>
            <integer>2</integer>
            <key>Hour</key>
            <integer>9</integer>
            <key>Minute</key>
            <integer>15</integer>
        </dict>
        <!-- Wednesday 9:15 AM -->
        <dict>
            <key>Weekday</key>
            <integer>3</integer>
            <key>Hour</key>
            <integer>9</integer>
            <key>Minute</key>
            <integer>15</integer>
        </dict>
        <!-- Thursday 9:15 AM -->
        <dict>
            <key>Weekday</key>
            <integer>4</integer>
            <key>Hour</key>
            <integer>9</integer>
            <key>Minute</key>
            <integer>15</integer>
        </dict>
        <!-- Friday 9:15 AM -->
        <dict>
            <key>Weekday</key>
            <integer>5</integer>
            <key>Hour</key>
            <integer>9</integer>
            <key>Minute</key>
            <integer>15</integer>
        </dict>
    </array>
    
    <key>StandardOutPath</key>
    <string>/Users/kuma/Desktop/Python/Stock_Bot/stockbot.log</string>
    
    <key>StandardErrorPath</key>
    <string>/Users/kuma/Desktop/Python/Stock_Bot/stockbot.error.log</string>
    
    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
    </dict>
    
    <key>RunAtLoad</key>
    <false/>
</dict>
</plist>

